# linux/pixhawk2xbee.service
[Unit]
Description=Pixhawk to XBee script.

[Service]
ExecStart=/usr/local/bin/pixhawk2xbee

[Install]
WantedBy=multi-user.target 
